# llm-gateway-console/deploy/llm-gateway-console.service
[Unit]
Description=LLM Gateway Console
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/llm-router-panel/llm-gateway-console/backend
EnvironmentFile=/opt/llm-router-panel/llm-gateway-console/backend/.env
ExecStart=/opt/llm-router-panel/llm-gateway-console/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8010
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
